# .config/systemd/user/firefox.service
[Unit]
Description=Firefox-nightly

[Service]
ExecStart=/usr/bin/firefox-nightly
Restart=always

[Install]
WantedBy=home.target
