# home/.config/systemd/user/firefox.service
[Unit]
Description=Firefox

[Service]
ExecStart=/usr/bin/firefox
Restart=always

[Install]
WantedBy=home.target
