# home/.config/systemd/user/home.target
[Unit]
Description=Xinitrc Stuff

[Install]
Alias=default.target
