# home/.config/systemd/user/ncmpcpp.service
[Unit]
Description=ncmpcpp

[Service]
ExecStart=/usr/bin/konsole --nofork -e ncmpcpp
Restart=always

[Install]
WantedBy=home.target
