# home/.config/systemd/user/udiskie.service
[Unit]
Description=devmon Service

[Service]
ExecStart=/usr/bin/udiskie -c /home/alk/.config/udiskie

[Install]
WantedBy=home.target

